# -config/autostart/remap-sculpt.desktop
[Desktop Entry]
Name=Remap Sculpt
Exec=remap-sculpt
Terminal=false
Type=Application
StartupNotify=false
